# templates/default/chef-rundeck.conf.erb
#
# This file is generated by Chef
#

# chef-rundeck bridge
# upstart service
start on runlevel [2345]
stop on shutdown
# automatically respawn
respawn
respawn limit 3 30
# generated by chef from attributes
exec su - <%= @bridge_user %> -s /bin/sh -c 'exec "$0" "$@"' -- "<%= @bridge_binary %>" <% @bridge_options.sort.each do |key, value| -%>--<%= key %> <%= value %> <% end %> 2>&1 > <%=@bridge_logdir%>/server.log
